# etc/NetworkManager/dispatcher.d/02default_route_workaround
#!/bin/sh -e

ip route flush proto boot
