# lib/systemd/system/alsa-state.service.d/disabled.conf
[Unit]
ConditionPathExists=!/lib/systemd/system/lxc-android-config.service
RefuseManualStart=true
