# etc/profile.d/disable-qtwebengine-gpu.sh
#!/bin/sh
# Disable the GPU for QtWebEngine on Android-based devices
# This file is part of lxc-android-config
export QT_WEBENGINE_DISABLE_GPU=1
