# lib/systemd/system/ofono.service.d/lxc-android-config.conf
[Unit]
Wants=rild.service
After=rild.service

[Service]
Environment=OFONO_RIL_DEVICE=ril
Environment=OFONO_RIL_NUM_SIM_SLOTS=1
ExecStart=
ExecStart=/usr/sbin/ofonod -P stktest,provision,sap,udev,dun,smart,hfp --nodetach
