# etc/init/enable-cpu-hotplugging.conf
description     "enable cpu hotplugging after Ubuntu Touch boot"

start on runlevel [2345]

task

script
sleep 60
for cpu in /sys/devices/system/cpu/cpu?/online; do
	chmod 644 $cpu
done
end script
