# etc/init/force-mtp.conf
start on android and started lightdm

task

script
        if ! android-gadget-service status mtp|grep -q enabled
        then
                android-gadget-service enable mtp
        fi
end script
