# etc/init/no-cpu-hotplug.conf
description     "disable cpu hotplugging during Ubuntu Touch boot"

start on local-filesystems

task

script
for cpu in /sys/devices/system/cpu/cpu?/online; do
        chmod 444 $cpu
done
end script
