# etc/init/ssh-keygen.conf
start on starting ssh

task

exec /usr/bin/ssh-keygen.sh